# gpuoff.service
[Unit]
Description=gpuoff

[Service]
ExecStart=/home/timothymiller/.startup/gpuoff.sh
